# sourcehut-man/sourcehut-man.service
[Unit]
Description=sourcehut-man

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/share/srht/man/run.py

