# sourcehut-git/sourcehut-git-api.service
[Unit]
Description=SourceHut git api

[Service]
Type=simple
WorkingDirectory=/usr/share/srht/git
ExecStart=/usr/bin/gitsrht-api

