# sourcehut-meta.service
[Unit]
Description=sourcehut-meta

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/share/sourcehut/meta.sr.ht/run.py
